# supervisord.conf
[supervisord]

nodaemon=true

[program:mysql]

command=/start.sh mysqld
